# letsautoencrypt.timer
[Unit]
Description=Auto sign https certs from Lets Encrypt

[Timer]
OnBootSec=10min
OnUnitActiveSec=2months

[Install]
WantedBy=timers.target

